# root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid redis redis-server /etc/redis/6379/redis-full.conf;
